# jupyter.conf
[program:jupyter]

directory=~/.jupyter
command=sudo jupyter-notebook
autostart=true
autorestart=true
stderr_logfile=/var/log/jupyter.err.log
stdout_logfile=/var/log/jupyter.out.log
startsecs = 0
